# aria2.service
[Unit]
Description=aria2

[Service]
Type=simple
User=yuhuang
ExecStart=$bin/aria2c --conf-path=$opt_prefix/aria2.conf --input-file=$opt_prefix/aria2.session
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=aria2.service
